# roles/system/templates/mnt-share.automount.j2
[Unit]
Description=automount {{ item.src }}
After=network.target multi-user.target

[Automount]
Where={{ item.path }}

[Install]
WantedBy=default.target
